# templates/default/supervisord.conf.erb
[program:<%=@app%>]
command=/bin/bash -c "service <%=@app%> start"